# templates/hermes-agent/hermes-dashboard.service
[Unit]
Description=Hermes Agent Web Dashboard
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment=PATH=%h/.hermes-venv/bin:%h/.local/bin:%h/.npm-global/bin:%h/bin:/usr/local/bin:/usr/bin:/bin
Environment=VENV=%h/.hermes-venv
ExecStart=%h/.hermes-venv/bin/hermes dashboard --host 0.0.0.0 --no-open --insecure
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
